# amavis/rootfs/etc/services.d/amavis/run
#!/usr/bin/with-contenv sh

echo "Starting Amavis daemon..."
s6-setuidgid amavis amavisd -c /etc/amavisd.conf
#amavisd -c /etc/amavisd.conf
