# combo/rootfs/etc/services.d/spamassassin/run
#!/usr/bin/with-contenv sh

sa-update

echo "Starting SpamAssassing daemon..."
exec spamd -s stderr
