# amavis/rootfs/etc/services.d/spamassassin/run
#!/usr/bin/with-contenv sh

sa-update -v

echo "Starting SpamAssassing daemon..."
exec spamd
