# combo/rootfs/etc/services.d/dovecot/run
#!/usr/bin/with-contenv sh

exec dovecot -F
