# clamsmtp/rootfs/etc/services.d/clamsmtp/run
#!/usr/bin/with-contenv sh

echo "Starting ClamSMTP daemon..."
#s6-setuidgid clamav clamsmtpd
exec clamsmtpd -d 4 2>&1
